# templates/django_supervisord.ini
[program:django_web_server]
command=/home/vagrant/venvs/webdev/bin/uwsgi --ini /vagrant/django18test/django18test/uwsgi.ini
user=vagrant
autorestart=true
stdout_logfile=/var/log/supervisor/jobs/django_app_server.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
redirect_stderr=true


[program:djangosocketio]
command=/home/vagrant/venvs/webdev/bin/python2.7 /vagrant/django18test/manage.py runserver_socketio --settings=django18test.settings.production
user=vagrant
autorestart=true
stdout_logfile=/var/log/supervisor/jobs/djangosocketio.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
redirect_stderr=true