# templates/django_supervisord.ini
[program:django_web_server]
command={{ venv_bin_path }}/uwsgi --ini {{ django_project_path }}/{{ django_project_name }}/uwsgi.ini
user={{ user_name }}
autorestart=true
stdout_logfile=/var/log/supervisor/jobs/django_server.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
redirect_stderr=true


[program:djangosocketio]
command={{ venv_bin_path }}/python2.7 {{ django_project_path }}/manage.py runserver_socketio --settings={{ django_project_name }}.settings.production
user={{ user_name }}
autorestart=true
stdout_logfile=/var/log/supervisor/jobs/django_socketio.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=5
stdout_capture_maxbytes=1MB
redirect_stderr=true